# server/sitxthsense.service
[Unit]
Description=Sense HAT API server

[Service]
Type=oneshot
ExecStart=/usr/bin/main.py

[Install]
WantedBy=multi-user.target

