# server/sixthsense.service
[Unit]
Description=Sense HAT API server
After=network.target

[Service]
ExecStart=/opt/sixthsense.py

[Install]
WantedBy=multi-user.target
